# packaging/arch/system-bridge.service
[Unit]
Description=VibeCast Autonomous System Bridge
After=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/bin/system-bridge
Restart=on-failure
RestartSec=5

# Privileged daemon for managing network, bluetooth, and system updates
User=vibecast

# Hardening options (can be adjusted later if they break zbus/Arch Toolkit)
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_TIME
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_TIME

[Install]
WantedBy=multi-user.target
